# assets/services/cortile.service
[Unit]
Description=window tiling manager
After=graphical.target

[Service]
ExecStart=/usr/local/bin/bsptile
Restart=always

[Install]
WantedBy=default.target
